# examples/com.example.hdhomerun-jellyfin-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.hdhomerun-jellyfin-refresh</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/YOUR_USER/Library/Application Support/hdhomerun-jellyfin/refresh_hdhomerun_jellyfin.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>43200</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>WorkingDirectory</key>
  <string>/Users/YOUR_USER/Library/Application Support/hdhomerun-jellyfin</string>
  <key>StandardOutPath</key>
  <string>/Users/YOUR_USER/Library/Logs/hdhomerun-jellyfin-refresh.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/YOUR_USER/Library/Logs/hdhomerun-jellyfin-refresh.err.log</string>
</dict>
</plist>

